# supervisord.conf
[program:notforsale_client]
command=/usr/local/bin/python -u /app/client.py